# src/main/resources/crontab
#Minutes(0-59) Hours(0-23) Day-of-Month(1-31) Month(1-12/JAN-DEC) Day-of-Week(0-6/SUN-SAT) javaMethod arg...
*/1 * * * * task.jcron.JCronTask1#main
42 16 * * * task.jcron.JCronTask2#run Hello World